# template/alsamixer-freeze.service
[Unit]
Description=Prevent applications to modify AlsaMixer settings

[Service]
ExecStart={{BINARY}}

[Install]
WantedBy=multi-user.target